# recipes-connectivity/bluetooth/files/bleadvhci.path
[Unit]
Description= RCP Beaconing 

[Path]
PathExists=/tmp/bluezup
Unit=bleadvhci.service

[Install]
WantedBy=multi-user.target
